# deploy/linode/systemd/jobpilot-api.service
[Unit]
Description=JobPilot FastAPI (uvicorn + Socket.IO)
After=network.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_USER__
WorkingDirectory=__APP_ROOT__/backend
EnvironmentFile=__APP_ROOT__/backend/.env
Environment=TECTONIC_PATH=__APP_ROOT__/backend/.bin/tectonic
ExecStart=__APP_ROOT__/backend/.venv/bin/uvicorn app.main:socket_app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
